# bash/network_reset/network_check_restart.service
[Unit]
Description=Virtual Distributed Ethernet

[Service]
ExecStart=bash /home/arnab/Scripts/network_reset.sh

[Install]
WantedBy=multi-user.target